# vps/etc/systemd/system/user@2000.service.d/10-pangolin.conf
# The pangolin user's manager starts the quadlet pod; make sure the config
# exists before it comes up.
[Unit]
After=pangolin-config-init.service
Requires=pangolin-config-init.service
