# vps/etc/systemd/system/user@2000.service.d/20-authentik.conf
# The quadlets reference authentik's podman secrets and /var/lib/authentik;
# both must exist before the user manager starts the services pod.
[Unit]
After=authentik-config-init.service
Requires=authentik-config-init.service
