# vps/etc/containers/systemd/users/2000/pangolin.container
[Container]
ContainerName=pangolin
Image=docker.io/fosrl/pangolin:ee-latest

Pod=services.pod

HealthCmd=["curl","-f","http://localhost:3001/api/v1/"]
HealthInterval=10s
HealthRetries=15
HealthTimeout=10s
Notify=healthy

Volume=/var/lib/pangolin/config:/app/config:z
Volume=/var/lib/pangolin/config/letsencrypt:/app/config/letsencrypt:ro,z

[Service]
Restart=always
RestartSec=5
# Cover slow image pulls and app startup; the 90s default start timeout
# marks the unit failed, which permanently kills units that depend on it.
TimeoutStartSec=900

[Install]
WantedBy=default.target
